# deploy/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm81 -F

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
